# cotabot-panel.service
[Unit]
Description=Cotabot Web Admin Panel
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/cotabot/COTABOT - DEV
ExecStart=/usr/bin/python3 /home/ubuntu/cotabot/COTABOT - DEV/web_admin/api.py
Restart=always
RestartSec=10

# Environment variables
Environment="PYTHONUNBUFFERED=1"
Environment="WEB_ADMIN_SECRET_KEY=cotabot-web-admin-secret-key-change-in-production"
Environment="WEB_ADMIN_API_KEY=cotabot-admin-2024"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cotabot-panel

[Install]
WantedBy=multi-user.target
